# dist/systemd/user/gmpd.service
[Unit]
Description = Play Music Proxy for MPD

[Service]
Type = simple
ExecStart = /usr/bin/node %h/.config/gmpd/server.js
WorkingDirectory = %h/.config/gmpd
Restart = always
RestartSec = 10

[Install]
WantedBy = multi-user.target
